# hello/hello.service
[Unit]
Description=Hello Service

[Service]
ExecStart=/usr/bin/python /demo/hello.py
