# web/etc/systemd/app.target
[Unit]
Description=App Target
Requires=app.service redis.service index.timer index.service