# web/etc/systemd/index.service
[Unit]
Description=Index DB

[Service]
ExecStart=/usr/bin/python /demo/web/index.py
