# web/etc/systemd/index.timer
[Unit]
Description=Index DB Timer

[Timer]
#OnCalendar=*:0/2
OnUnitActiveSec=5sec
AccuracySec=1us
