# web/etc/systemd/redis.service
[Unit]
Description=Redis

[Service]
ExecStart=/usr/bin/redis-server --port 7777
